# notifier/crontab
0 9 * * * /notify.sh
0 21 * * * /notify.sh
@reboot /notify.sh
